# python-services/embedding-service.service
[Unit]
Description=WhatsApp Auto-Reply Embedding Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/mensajeria-app/python-services
Environment="PATH=/home/ubuntu/mensajeria-app/python-services/venv/bin:/usr/bin:/usr/local/bin"
ExecStart=/home/ubuntu/mensajeria-app/python-services/venv/bin/python embedding_service.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=embedding-service

[Install]
WantedBy=multi-user.target
